# script/contentful.crontab
* * * * *  cd /opt/www.aptible.com && bundle exec rake contentful:pull
